# beeppp-server.socket
[Unit]
Description=Beeppp socket

[Socket]
ListenStream=/run/beeppp-server.sock

[Install]
WantedBy=sockets.target